# docker/supervisor/laravel_schedule.conf
[program:laravel_schedule]
command=bash -c "sleep 60 && php /opt/artisan schedule:run"
autostart=true
autorestart=true
stderr_logfile=/var/log/long.err.log
stdout_logfile=/var/log/long.out.log